# print-server.service
[Unit]
Description=Thermal Print Server
After=network.target

[Service]
Type=simple
User=berend
WorkingDirectory=/home/berend/thermal-printer
ExecStart=/home/berend/.nvm/versions/node/v22.22.0/bin/node --experimental-strip-types src/index.ts
Restart=always
RestartSec=3
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
